# roles/asdf/templates/etc/profile.d/asdf.sh.j2
{% if os == "fedora" -%}
ASDF_BIN="{{ user_home_dir }}/.local/bin"

if ! [[ $PATH =~ $ASDF_BIN ]]; then
  PATH="$PATH:$ASDF_BIN"
fi
{% endif %}
