# roles/asdf/templates/etc/profile.d/asdf.sh.j2
{% if common_os | lower == "fedora" -%}
ASDF_BIN="{{ common_user_home_dir }}/.local/bin"

if ! [[ $PATH =~ $ASDF_BIN ]]; then
  PATH="$PATH:$ASDF_BIN"
fi
{% endif %}
